# deploy/systemd/banall-bot.service.template
[Unit]
Description=BanAll Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/scripts/run_bot.sh
Restart=always
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGINT

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_DIR__/logs __APP_DIR__/bot.log

[Install]
WantedBy=multi-user.target

